# rituals/netmon.timer
[Unit]
Description=Network Monitor Timer - Run speedtest every 30 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
